# lsfmon/install/accounting.cron.tmpl
PATH=/sbin:/bin:/usr/sbin:/usr/bin
05 */1 * * * root BASE_DIR/bin/launch_accounting.sh >> /tmp/cron_accounting.log 2>&1
