# lsfmon/install/overview.cron.tmpl
PATH=/sbin:/bin:/usr/sbin:/usr/bin
*/3 * * * * root BASE_DIR/bin/launch_overview.sh >> /tmp/cron_overview.log 2>&1
