# jobview-condor/cron/overview.cron
PATH=/sbin:/bin:/usr/sbin:/usr/bin
MAILTO=subir.sarkar@cern.ch
00-59/3 * * * * root /opt/jobview_t2/bin/launch_overview.sh >> /tmp/cron_overview.log 2>&1
