# jobmon-glidein-frontend/install/voUserList.cron.tmpl
PATH=/sbin:/bin:/usr/sbin:/usr/bin
05 1,13 * * * root BASEDIR/jobmon/bin/launch_voUserList.sh > BASEDIR/jobmon/log/cron_voUserList.log 2>&1
