# jobview-pbs/bin/overview.cron
PATH=/sbin:/bin:/usr/sbin:/usr/bin
MAILTO=subir.sarkar@pi.infn.it
00-59/3 * * * * root /opt/jobview/bin/launch_overview.sh > /tmp/cron_overview.log 2>&1
